# conf/supervisor/conf.d/api.conf
[program:puma]
command=bundle exec puma -C config/puma.rb
directory=/api
stdout_logfile=/var/log/supervisor/puma-access.log
stderr_logfile=/var/log/supervisor/puma-error.log
autostart=true
autorestart=true

[program:sidekiq]
command=bundle exec sidekiq -C config/sidekiq.yml | tee ./log/sidekiq.log
directory=/api
stdout_logfile=/var/log/supervisor/sidekiq-access.log
stderr_logfile=/var/log/supervisor/sidekiq-error.log
autostart=true
autorestart=true